# deploy/hospital_system.service
# systemd service file - keeps the app running 24/7, auto-restarts on crash/reboot
# Install: sudo cp deploy/hospital_system.service /etc/systemd/system/
# Then:    sudo systemctl daemon-reload && sudo systemctl enable --now hospital_system

[Unit]
Description=Hospital Management System (Gunicorn)
After=network.target postgresql.service

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/hospital_system
EnvironmentFile=/var/www/hospital_system/.env
ExecStart=/var/www/hospital_system/venv/bin/gunicorn -c deploy/gunicorn_config.py hospital_system.wsgi
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
